# files/etc/uci-defaults/10_disable_services
/etc/init.d/tailscale disable
echo "#20 5 2 * * /scripts/upgrade.sh online needback stable" >> /etc/crontabs/root

exit 0