# files/etc/uci-defaults/10_disable_services
/etc/init.d/tailscale disable
echo "#20 5 * * 1 /scripts/upgrade.sh online needback" >> /etc/crontabs/root

exit 0